# files/pimk.service
[Unit]
Description = Pimk USB keyboard gadget

[Service]
Type=simple
ExecStart = /usr/local/bin/pimk

[Install]
WantedBy = multi-user.target
